# .factory/worktrees/task-bee4f48c/service/systemd/claude-code-model-gateway.service
[Unit]
Description=Claude Code Model Gateway Service
Documentation=https://github.com/anthropic/claude-code-model-gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=claude-gateway
Group=claude-gateway

# Environment and configuration
EnvironmentFile=-/etc/claude-code-model-gateway/environment
WorkingDirectory=/opt/claude-code-model-gateway

# Service executable — uses the installed CLI entry point
ExecStartPre=/opt/claude-code-model-gateway/venv/bin/claude-code-model-gateway version
ExecStart=/opt/claude-code-model-gateway/venv/bin/claude-code-model-gateway gateway \
    --host ${GATEWAY_HOST:-0.0.0.0} \
    --port ${GATEWAY_PORT:-8080} \
    --timeout ${GATEWAY_TIMEOUT:-300} \
    --log-format ${GATEWAY_LOG_FORMAT:-json} \
    --log-file ${GATEWAY_LOG_FILE:-/var/log/claude-code-model-gateway/gateway.log}
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Timeouts
TimeoutStartSec=30
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=no
SystemCallArchitectures=native

# File access
ReadWritePaths=/var/log/claude-code-model-gateway
ReadWritePaths=/var/lib/claude-code-model-gateway
ReadOnlyPaths=/etc/claude-code-model-gateway
ReadOnlyPaths=/opt/claude-code-model-gateway

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-code-model-gateway

[Install]
WantedBy=multi-user.target
